# Dublin/bulk/weekly/crontab
0 14 * * Tue /usr/local/bin/python /root/main.py
