# Dublin/bulk/daily/crontab
5 20 * * * /usr/local/bin/python /root/main.py
