# Dublin/crontab/crontab
*/20 * * * * /home/matlab/croco-op/Dublin/alerts/checklogs.sh
30 16 * *  Tue /home/matlab/croco-op/Dublin/run/HC/hindcast.sh
30 21 * * * /home/matlab/croco-op/Dublin/run/FC/forecast.sh