# Dublin/input/HC/crontab
20 15 * * Tue /usr/local/bin/python /root/main.py
